# src/robot_control.service
[Unit]
Description=Robot Web Control Service
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/barry/mcar/src
ExecStart=/bin/bash /home/barry/mcar/src/start_robot_control.sh
Restart=on-failure
RestartSec=5
StartLimitIntervalSec=60
StartLimitBurst=3

[Install]
WantedBy=multi-user.target 